# supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins]
user=jenkins
autostart=true
autorestart=true
command=/usr/local/bin/jenkins.sh
redirect_stderr=true
# stdout_logfile=/var/log/jenkins/%(program_name)s.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1
environment = JENKINS_HOME="/var/jenkins_home",HOME="/var/jenkins_home",USER="jenkins",COPY_REFERENCE_FILE_LOG="/var/log/copy_reference_file.log"

# [program:jenkins]
# command=/usr/local/bin/jenkins.sh
# user=jenkins
# priority=2
# autostart=true
# autorestart=unexpected
# stdout_logfile=/tmp/jenkins_out.log
# stderr_logfile=/tmp/jenkins_err.log

[program:ibmdbb]
command=/usr/local/bin/start.sh
user=root
priority=2
autostart=true
autorestart=unexpected
stdout_logfile=/tmp/ibmdbb_out.log
stderr_logfile=/tmp/ibmdbb_err.log